# trunk/spec_files/ext-sources/saslauthd.xml
<? xml version ='1 .0 '?>
<! DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='saslauthd'>

<service
	name='application/security/saslauthd'
	type='service'
	version='1?>

	<create_default_instance enabled='false' />
	<single_instance/> <single_instance/>

	<dependency
		name='filesystem'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/system/filesystem/local' />
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/usr/gnu/sbin/saslauthd -a pam'
		timeout_seconds='30' >
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec=':kill'
		timeout_seconds='60' />

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>Cyrus SASL server</loctext>
		</common_name>
		<documentation>
			<manpage title='saslauthd' section='8' manpath='/usr/gnu/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
