# misc/belenix-branding/profile.d/kde.pagein
/usr/X11/lib/libX11.so.6
/usr/sfw/lib/libfreetype.so.6
/usr/foss/lib/libkdecore.so
/usr/foss/lib/libDCOP.so
/usr/foss/lib/libart_lgpl_2.so
/usr/foss/lib/libidn.so
/usr/foss/lib/libkdefx.so
/usr/foss/lib/libqt-mt.so
/usr/foss/lib/libstdc++.so
/usr/foss/lib/libksplashthemes.so
/usr/foss/lib/libkio.so
/usr/foss/lib/libkdeui.so
