# misc/belenix-branding/profile.d/xfce.pagein
/usr/X11/lib/libX11.so.6
/usr/foss/lib/libxfcegui4.so
/usr/foss/lib/libgtk-x11-2.0.so
/usr/foss/lib/libgdk-x11-2.0.so
/usr/foss/lib/libpango-1.0.so
/usr/foss/lib/libgobject-2.0.so
/usr/foss/lib/libglib-2.0.so
/usr/foss/lib/libcairo.so
