# branches/belenix.0.7/misc/belenix-branding/profile.d/kde.minimal.pagein
/usr/X11/lib/libX11.so.6
/usr/lib/libkdecore.so
/usr/lib/libqt-mt.so
