# misc/belenix-branding2/profile.d/env.sh
PATH=/usr/gnu/bin:/sbin:/usr/sbin:/usr/bin:/usr/ccs/bin:/usr/X11/bin:/usr/openwin/bin:/opt/DTT/bin
export PATH

MANPATH=/usr/man:/usr/share/man:/usr/gnu/share/man:/usr/X11/man:/opt/DTT/Man:/usr/sfw/share/man
export MANPATH

PKG_CONFIG_PATH=/usr/lib/pkgconfig:/usr/X11/lib/pkgconfig:/usr/gnu/lib/pkgconfig:/usr/sfw/lib/pkgconfig
export PKG_CONFIG_PATH

TMPDIR=/tmp
export TMPDIR

#KDEVARTMP=/tmp
#export KDEVARTMP

G_BROKEN_FILENAMES=yes
export G_BROKEN_FILENAMES

QTDIR=/usr
export QTDIR

QMAKESPEC=${QTDIR}/mkspecs/solaris-g++
export QMAKESPEC

SASL_PATH=/usr/gnu/lib/sasl2
export SASL_PATH

SSH_AGENT_EXE=/usr/bin/ssh-agent
export SSH_AGENT_EXE

XFCE4_SESSION_EXE=/usr/bin/xfce4-session
export XFCE4_SESSION_EXE

XDG_DATA_DIRS=/usr/share
export XDG_DATA_DIRS
